# deploy/martin.service
[Unit]
Description=Martin Vector Tile Server (PostGIS)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/martin --listen-addresses 127.0.0.1:3000 postgresql://ubuntu@localhost/buoy_mappers
Restart=on-failure
RestartSec=5s
StandardOutput=append:/var/log/martin.log
StandardError=append:/var/log/martin.error.log
SyslogIdentifier=martin

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
